# dist/supervisord.conf
[supervisord]
nodaemon=true

[program:cowrie]
user=tpot
directory=/opt/cowrie/
command=bash -c 'cp /opt/cowrie/misc/userdb.txt /data/cowrie/misc/ && exec /usr/local/bin/twistd --nodaemon -y cowrie.tac -l /data/cowrie/log/cowrie.log --pidfile /var/run/cowrie/cowrie.pid cowrie'
autorestart=true

